# contrib/systemd/pfxhttp.service
[Unit]
Description=PfxHTTP Postfix-to-HTTP server
After=network.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/etc/default/pfxhttp
# Be explicit: point to the config file in /etc to avoid search path surprises
ExecStart=/usr/local/sbin/pfxhttp --config /etc/pfxhttp/pfxhttp.yml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2s

# Notes:
# - The service starts as root to bind privileged ports or create Unix sockets
#   in protected directories (e.g. Postfix spool). The application then drops
#   privileges internally to the user/group configured in pfxhttp.yml.
# - We deliberately avoid over-hardening here to keep /etc visible and allow
#   setuid/setgid. Add hardening later only after verifying functionality.

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pfxhttp
MemoryMax=50M
CPUQuota=10%

# Minimal sandboxing that doesn't interfere with config access or uid/gid changes
PrivateTmp=true
# If you enable ProtectSystem= or RestrictSUIDSGID= in the future, verify that
# the config can still be read and the internal privilege drop still works.
# If you need to create a Unix socket in Postfix' private spool, you may need:
# ReadWritePaths=/var/spool/postfix-submission/private

LimitNOFILE=1024

[Install]
WantedBy=multi-user.target